# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:gunicorn]
command=/home/caleb/venv/bin/gunicorn -w 1 blackshakara.project.wsgi:application --bind 0.0.0.0:8000
directory=/home/caleb/backend/
user=caleb
autostart=true
autorestart=true
redirect_stderr=true
stderr_logfile=/var/log/gunicorn.err.log
stdout_logfile=/var/log/gunicorn.out.log

[program:celery]
command=/home/caleb/venv/bin/celery -A blackshakara.project worker -l info --pool=solo
directory=/home/caleb/backend/
user=caleb
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/celery.err.log
stdout_logfile=/var/log/celery.out.log

[program:celery_beat]
command=/home/caleb/venv/bin/celery -A blackshakara.project beat -l INFO
directory=/home/caleb/backend/
user=caleb
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/celery_beat.err.log
stdout_logfile=/var/log/celery_beat.out.log
